# supervisord.conf
[supervisord]
nodaemon=true

[program:alerta]
command=/bin/bash -c "alertad run --host 0.0.0.0 --port 5000"

[program:alerta-mailer]
command=/bin/bash -c "alerta-mailer"
